# desktop/.config/systemd/user/dlcleaner.service
[Unit]
Description=Backup of my apache website

[Service]
Type=oneshot
ExecStart=/home/prashant/.hscripts/cleaner
