# thinkpad/.config/systemd/user/dlcleaner.service
[Unit]
Description=Cleans downloads

[Service]
Type=oneshot
ExecStart=/bin/bash /home/prash/.scripts/dlcleaner.sh
