# thinkpad/.config/systemd/user/dlcleaner.service
[Unit]
Description=Cleans downloads

[Service]
Type=oneshot
ExecStart=/home/prash/.hscripts/cleaner
