# mba/service/iw/run
#!/bin/bash
iwctl station wlp3s0 connect "modem1"
