# newdesktop/.config/systemd/user/dlcleaner.service
[Unit]
Description=Cleans downloads folder

[Service]
Type=oneshot
ExecStart=/bin/bash /home/prashant/.scripts/dlcleaner.sh
